# price-monitor.service
[Unit]
Description=综合秘书机器人 Secretary Bot v2
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu                                          # 更改为你的用户名
Group=ubuntu                                         # 更改为你的用户组
WorkingDirectory=/home/ubuntu/PriceMonitor           # 更改为项目实际路径
EnvironmentFile=/home/ubuntu/PriceMonitor/.env       # 从 .env 加载所有密钥

ExecStart=/usr/bin/python3 /home/ubuntu/PriceMonitor/bot/app.py

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/PriceMonitor
ProtectHome=read-only

[Install]
WantedBy=multi-user.target