# price-monitor.service
[Unit]
Description=USD/CNH and Brent Oil Price Monitor Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu  # 更改为你的用户名
Group=ubuntu # 更改为你的用户组
WorkingDirectory=/home/ubuntu/price-monitor  # 更改为项目目录
Environment="PATH=/home/ubuntu/price-monitor/venv/bin"  # 如果使用虚拟环境
Environment="ALPHA_VANTAGE_API_KEY=your_api_key_here"  # 可选：使用环境变量
Environment="TELEGRAM_BOT_TOKEN=your_bot_token_here"
Environment="TELEGRAM_CHAT_ID=your_chat_id_here"

# 使用虚拟环境中的Python
ExecStart=/home/ubuntu/price-monitor/venv/bin/python /home/ubuntu/price-monitor/price_monitor.py

# 或者直接使用系统Python（如果不使用虚拟环境）
# ExecStart=/usr/bin/python3 /home/ubuntu/price-monitor/price_monitor.py

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/price-monitor
ProtectHome=read-only

[Install]
WantedBy=multi-user.target